# loadboard.service
[Unit]
Description=Loadboard FastAPI (Uvicorn)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=www-data
WorkingDirectory=/var/www/loadboard-fastapi
ExecStart=/var/www/loadboard-fastapi/venv/bin/uvicorn main:app --uds /run/loadboard-fastapi.sock --workers 4 --loop uvloop
Restart=always
RestartSec=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
